# qtrigdoppler-remote.service
[Unit]
Description=QTrigdoppler Remote Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/path/to/QTrigdoppler
ExecStart=/usr/bin/node remote_server.js
Restart=on-failure
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=qtrigdoppler-remote

[Install]
WantedBy=multi-user.target
